# wekan-mcp.service
[Unit]
Description=Wekan MCP Server
After=network.target

[Service]
Type=simple
User=wekan-mcp
Group=wekan-mcp
WorkingDirectory=/opt/wekan-mcp
ExecStart=/opt/wekan-mcp/venv/bin/python /opt/wekan-mcp/server.py
Restart=always
RestartSec=5
Environment=PATH=/opt/wekan-mcp/venv/bin:/usr/bin:/bin

# Security hardening (Issue #15)
ProtectSystem=strict
ReadWritePaths=/opt/wekan-mcp
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
